# fps-arena.service
[Unit]
Description=FPS Arena (枪战突击 COMBAT STRIKE)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/fps-arena
ExecStart=/usr/bin/node /opt/fps-arena/server/index.js
Restart=on-failure
RestartSec=3
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=DB_PATH=/opt/fps-arena/data/app.sqlite

[Install]
WantedBy=multi-user.target
